# archived/v1/share/dar-inc-backup.timer
[Unit]
Description=dar-backup INC timer

[Timer]
OnCalendar=*-*-04/3 19:03:00

[Install]
WantedBy=timers.target


