# archived/v1/share/dar-cleanup.timer
[Unit]
Description=dar-cleanup DIFF & INC timer

[Timer]
OnCalendar=*-*-* 21:07:00

[Install]
WantedBy=timers.target


