# templates/systemd/dar-inc-backup.service.template
[Unit]
Description=dar-backup INC
StartLimitIntervalSec=120
StartLimitBurst=1
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart="@@DAR_BACKUP_DIR@@/bin/dar-inc-backup.sh"
