# share/alert-upcoming-full-backup.timer
[Unit]
Description=alert upcoming FULL dar-backup

[Timer]
OnCalendar=*-12-28 10:03:00

[Install]
WantedBy=timers.target


