# share/dar-backup.timer
[Unit]
Description=dar backup timer

[Timer]
OnCalendar=Sat *-*-* 13:31:30
Unit=dar-backup.service

[Install]
WantedBy=timers.target


