# share/dar-inc-backup.timer
[Unit]
Description=dar-backup INC timer

[Timer]
OnCalendar=*-02..12-04/3 19:03:00

[Install]
WantedBy=timers.target


